# rpi-lcd-monitor.service
[Unit]
Description=Raspberry Pi LCD System Monitor and VLC Controller
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/.venv/bin/python3 __APP_DIR__/src/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

